# build/etc/supervisor.d/nodervisor.ini
[program:nodervisor]
command=/usr/bin/npm start
directory=/usr/local/nodervisor
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
startretries=3
stopsignal=QUIT
stopwaitsecs=10